# supervisord.conf
[supervisord]
nodaemon = true

[program:sshd]
command = /usr/sbin/sshd -D

[inet_http_server]
port = *:9001
